# infra/picanm/systemd/n2k-raw-spool-health.service
[Unit]
Description=Check picanm raw N2K spool disk health

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check-raw-spool-space
