# infra/pi5nvme/systemd/boat-raw-n2k-import.timer
[Unit]
Description=Import mirrored raw NMEA 2000 logs periodically
Documentation=file:/home/jack/boat-data-platform/docs/2026-07-03-pi5nvme-incident-and-picanm-status.md
# Do not enable this timer on a live host without an approved import window.
# The paired service is also gated by /etc/boat-data-platform/allow-raw-n2k-import.

[Timer]
OnBootSec=30min
OnUnitActiveSec=6h
Persistent=false
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
