# infra/pi5nvme/systemd/boat-signalk-collector.service
[Unit]
Description=Collect Signal K deltas into TimescaleDB
Wants=network-online.target postgresql.service signalk-pi5nvme.service
After=network-online.target postgresql.service signalk-pi5nvme.service

[Service]
Type=simple
User=jack
Group=jack
WorkingDirectory=/home/jack/boat-data-platform
Environment=SIGNALK_WS_URL=ws://127.0.0.1:3001/signalk/v1/stream?subscribe=none
EnvironmentFile=/etc/boat-data-platform/db.env
ExecStart=/usr/bin/node scripts/collect-signalk-to-timescale.mjs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
